# systemd/systemd/user/dms-session.target
[Unit]
Description=DMS (Dank Material Shell) Session Target
PartOf=graphical-session.target
After=graphical-session.target
