# systemd/systemd/user/niri-session.target
[Unit]
Description=Niri Session Target
BindsTo=niri.service
Requires=graphical-session.target
After=graphical-session.target
After=niri.service
