# systemd/systemd/user/dms.service.d/session.conf
[Unit]
# Scope DMS to the DMS session target, not all graphical sessions
PartOf=
PartOf=dms-session.target

[Service]
# Only restart on crashes (SIGABRT etc.), not clean stops from session teardown
Restart=on-failure
