# systemd/systemd/user/plasma-xdg-desktop-portal-kde.service.d/no-plasma-core-ordering.conf
[Unit]
# Remove After=plasma-core.target to prevent D-Bus activation from cascading
# into 26 waiting jobs that deadlock the user systemd instance.
# The portal only needs kwin_wayland to be running, not the full Plasma core.
After=
