# systemd/systemd/user/gamescope-session.service.d/override.conf
[Service]
UnsetEnvironment=DISPLAY XAUTHORITY WAYLAND_DISPLAY
